# roles/postgresql/templates/postgresql.logrotate.service.j2
[Unit]
  Description=postgresql.logrotate

[Service]
  ExecStart=/bin/bash {{ pgdata_dir }}log/postgresql.logrotate.sh
  Type=oneshot
  User=postgres
  WorkingDirectory={{ pgdata_dir }}
