# roles/postgresql/templates/postgresql.logrotate.timer.j2
[Unit]
  Description=postgresql.logrotate

[Timer]
  OnCalendar={{ logrotate.on_calendar }}
